# ps3/package/etc/udev/rules.d/21-ps3fs.rules
KERNEL=="spufs", ACTION=="add", RUN+="mkdir -p /spu", RUN+="mount -t spufs spu /spu"
