# ps3/package/etc/init/ps3vramswap.conf
# ps3vramswap - Mount PS3 VRAM swap on boot
#
# This helper mounts filesystems in the correct order as the devices
# and mountpoints become available.

description	"Mount PS3 VRAM swap on boot"

start on startup
stop on starting rcS

task
exec mkswap /dev/ps3vram
