# ps3/package/etc/init/setdefaultrunlevel.conf
description "sets the default runlevel"

start on filesystem

console output

task

script
    lvl=2 
    telinit $lvl
end script
